# supervisord.conf
; =============================================================================
; Supervisord — Metasploitable3 ARM64 (extended CVE coverage May 2026)
; =============================================================================

[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
user=root

[unix_http_server]
file=/var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

; SSH — CVE-2024-6387 regreSSHion
[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/sshd.out.log
stderr_logfile=/var/log/supervisor/sshd.err.log

; vsftpd — CVE-2011-2523 (emulated)
[program:vsftpd]
command=/usr/sbin/vsftpd /etc/vsftpd.conf
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/vsftpd.out.log
stderr_logfile=/var/log/supervisor/vsftpd.err.log

; Backdoor sim port 6200 (CVE-2011-2523 equivalent)
[program:backdoor-sim]
command=/usr/local/bin/backdoor-sim
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/backdoor.out.log
stderr_logfile=/var/log/supervisor/backdoor.err.log

; ProFTPd — CVE-2015-3306 mod_copy RCE
[program:proftpd]
command=/usr/sbin/proftpd --nodaemon --config /etc/proftpd/proftpd.conf
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/proftpd.out.log
stderr_logfile=/var/log/supervisor/proftpd.err.log

; Apache — CVE-2024-38474 + CVE-2012-1823 + CVE-2014-6271 Shellshock
[program:apache2]
command=/usr/sbin/apache2ctl -D FOREGROUND
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/apache2.out.log
stderr_logfile=/var/log/supervisor/apache2.err.log
environment=APACHE_RUN_USER="www-data",APACHE_RUN_GROUP="www-data",APACHE_PID_FILE="/var/run/apache2/apache2.pid",APACHE_RUN_DIR="/var/run/apache2",APACHE_LOCK_DIR="/var/lock/apache2",APACHE_LOG_DIR="/var/log/apache2"

; MariaDB — CVE-2021-27928 + mysql_udf_payload
[program:mariadb]
command=/usr/bin/mysqld_safe
autostart=true
autorestart=true
user=mysql
stdout_logfile=/var/log/supervisor/mariadb.out.log
stderr_logfile=/var/log/supervisor/mariadb.err.log

; Tomcat — CVE-2025-24813 + CVE-2020-1938 Ghostcat + CVE-2021-44228 Log4Shell
[program:tomcat]
command=/usr/share/tomcat9/bin/catalina.sh run
autostart=true
autorestart=true
user=tomcat
environment=JAVA_HOME="/usr/lib/jvm/default-java",CATALINA_HOME="/usr/share/tomcat9"
stdout_logfile=/var/log/supervisor/tomcat.out.log
stderr_logfile=/var/log/supervisor/tomcat.err.log

; Samba — CVE-2017-7494 SambaCry + CVE-2007-2447
[program:smbd]
command=/usr/sbin/smbd --foreground --no-process-group
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/smbd.out.log
stderr_logfile=/var/log/supervisor/smbd.err.log

[program:nmbd]
command=/usr/sbin/nmbd --foreground --no-process-group
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/nmbd.out.log
stderr_logfile=/var/log/supervisor/nmbd.err.log

; Redis — CVE-2022-0543 Lua sandbox escape + unauthenticated RCE
[program:redis]
command=/usr/bin/redis-server /etc/redis/redis.conf
autostart=true
autorestart=true
user=redis
stdout_logfile=/var/log/supervisor/redis.out.log
stderr_logfile=/var/log/supervisor/redis.err.log

; UnrealIRCd — CVE-2010-2075 backdoor (build from source or simulation)
[program:ircd]
command=bash -c '/usr/local/sbin/ircd 2>/dev/null || /usr/local/bin/ircd-sim'
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/ircd.out.log
stderr_logfile=/var/log/supervisor/ircd.err.log
