# supervisord.conf
; =============================================================================
; Supervisord — Metasploitable3 ARM64 Service Manager
; =============================================================================

[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
user=root

[unix_http_server]
file=/var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/sshd.err.log
stdout_logfile=/var/log/supervisor/sshd.out.log

[program:vsftpd]
command=/usr/sbin/vsftpd /etc/vsftpd.conf
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/vsftpd.err.log
stdout_logfile=/var/log/supervisor/vsftpd.out.log

[program:backdoor-sim]
command=/usr/local/bin/backdoor-sim
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/backdoor.err.log
stdout_logfile=/var/log/supervisor/backdoor.out.log

[program:apache2]
command=/usr/sbin/apache2ctl -D FOREGROUND
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/apache2.err.log
stdout_logfile=/var/log/supervisor/apache2.out.log
environment=APACHE_RUN_USER="www-data",APACHE_RUN_GROUP="www-data",APACHE_PID_FILE="/var/run/apache2/apache2.pid",APACHE_RUN_DIR="/var/run/apache2",APACHE_LOCK_DIR="/var/lock/apache2",APACHE_LOG_DIR="/var/log/apache2"

[program:mariadb]
command=/usr/bin/mysqld_safe
autostart=true
autorestart=true
user=mysql
stderr_logfile=/var/log/supervisor/mariadb.err.log
stdout_logfile=/var/log/supervisor/mariadb.out.log

[program:tomcat]
command=/usr/share/tomcat9/bin/catalina.sh run
autostart=true
autorestart=true
user=tomcat
environment=JAVA_HOME="/usr/lib/jvm/default-java",CATALINA_HOME="/usr/share/tomcat9"
stderr_logfile=/var/log/supervisor/tomcat.err.log
stdout_logfile=/var/log/supervisor/tomcat.out.log

[program:smbd]
command=/usr/sbin/smbd --foreground --no-process-group
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/smbd.err.log
stdout_logfile=/var/log/supervisor/smbd.out.log

[program:nmbd]
command=/usr/sbin/nmbd --foreground --no-process-group
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/nmbd.err.log
stdout_logfile=/var/log/supervisor/nmbd.out.log
